# pkg/lpp-custom-application/DEBIAN/postinst
#!/bin/bash

systemctl daemon-reload
systemctl enable lpp-custom-application
systemctl restart lpp-custom-application
