# pkg/lpp-custom-application/DEBIAN/prerm
#!/bin/bash

systemctl daemon-reload
systemctl disable lpp-custom-application
systemctl stop lpp-custom-application
